# flyprint.service
[Unit]
Description=Fly Print Edge Service
Documentation=https://github.com/ecnu/fly-print
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ecnu
Group=ecnu
WorkingDirectory=/opt/flyprint
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/flyprint/.venv/bin/python /opt/flyprint/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
